# project_aware.cron
0 0 * * * /usr/local/bin/php /your/path/ProjectAware/index.php --uri=gather/index/daily
0 0,12 * * * /usr/local/bin/php /your/path/ProjectAware/index.php --uri=gather/index/twice_daily
0 0 * * 0 /usr/bin/foo /usr/local/bin/php /your/path/ProjectAware/index.php --uri=gather/index/weekly
0 0 * * 0,4 /usr/bin/foo /usr/local/bin/php /your/path/ProjectAware/index.php --uri=gather/index/twice_weekly
0 * 1 * * /usr/bin/foo/usr/bin/foo /usr/local/bin/php /your/path/ProjectAware/index.php --uri=gather/index/monthly
0 * 1,15 * * /usr/bin/foo /usr/local/bin/php /your/path/ProjectAware/index.php --uri=gather/index/twice_monthly
